# src/alexa@.service
[Unit]
Description=FAUST CTF alexa service

[Service]
ExecStart=python manage.py runserver 0.0.0.0:8888
WorkingDirectory=/srv/alexa
User=alexa
